# deploy/hetzner/systemd/chubby-mcp.service
# systemd service for Chubby MCP Server (Node.js)
#
# Install:
#   cp deploy/hetzner/systemd/chubby-mcp.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now chubby-mcp

[Unit]
Description=Chubby MCP Server (Node.js)
After=network-online.target chubby-sim.service
Wants=network-online.target
Requires=chubby-sim.service

[Service]
Type=simple
User=chubby
Group=chubby
WorkingDirectory=/opt/chubby/apps/mcp-server
ExecStart=/usr/bin/node dist/server.js
Environment=PORT=8080
Environment=NODE_ENV=production
Environment=SIMULATION_SERVICE_URL=http://localhost:3002

Restart=always
RestartSec=3

# Resource limits
MemoryMax=512M
MemoryHigh=384M
CPUQuota=200%
LimitNOFILE=4096

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/chubby

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chubby-mcp

[Install]
WantedBy=multi-user.target
