# deploy/hetzner/systemd/chubby-mcp.service
# systemd service for Chubby MCP Server (Go binary)
#
# Install:
#   cp deploy/hetzner/systemd/chubby-mcp.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now chubby-mcp

[Unit]
Description=Chubby MCP Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=chubby
Group=chubby
WorkingDirectory=/opt/chubby

# --- Go binary (default, recommended) ---
ExecStart=/opt/chubby/server
Environment=PORT=8080

# --- Node.js alternative (uncomment if using Node.js server) ---
# ExecStart=/usr/bin/node /opt/chubby/node/dist/server.js
# Environment=PORT=8080
# Environment=NODE_ENV=production
# Environment=WIDGET_VIEWER_URL=https://widget.chubby.fyi

Restart=always
RestartSec=3

# Resource limits: hard cap to prevent runaway costs
MemoryMax=512M
MemoryHigh=384M
CPUQuota=200%
LimitNOFILE=4096

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/opt/chubby

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chubby-mcp

[Install]
WantedBy=multi-user.target
